# ufodeck.service
[Unit]
Description=ufodeck
After=libvirtd.socket

[Service]
Type=simple
ExecStart=/usr/bin/ufodeck start
